# raspberry_pi/camera_streaming.service
[Unit]
Description=Camera Streaming Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/rover_cameras_base/raspberry_pi
ExecStart=/usr/bin/python3 /home/pi/rover_cameras_base/raspberry_pi/pi_streaming_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target